# scripts/launchd/com.yordam.agent.reorg.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.yordam.agent.reorg</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/yordamkocatepe/bin/yordam-agent</string>
      <string>reorg</string>
      <string>/path/to/folder</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/yordam-agent.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/yordam-agent.err</string>
  </dict>
</plist>
